# ops/launchd/com.seanburdges.lariat.archive-stale.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seanburdges.lariat.archive-stale</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/seanburdges/Dev/Lariat &amp;&amp; npm run archive:stale</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>WorkingDirectory</key>
  <string>/Users/seanburdges/Dev/Lariat</string>
  <key>StandardOutPath</key>
  <string>/Users/seanburdges/Dev/Lariat/logs/archive-stale.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/seanburdges/Dev/Lariat/logs/archive-stale.err.log</string>
</dict>
</plist>
